# rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

# Wait a few seconds for network connections to stabilize
echo "rc.local: Waiting a few seconds for network connections to stabilize..."
sleep 5

# Print information about network interfaces
echo "--------------------------------------------------------------------------------"
echo "Information about network interfaces:"
ifconfig

# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  echo "--------------------------------------------------------------------------------"
  printf "My IP address is %s\n" "$_IP"
fi

# Print information about network connections
echo "--------------------------------------------------------------------------------"
echo "Information about network connections:"
iwconfig

# Start light control script
echo "--------------------------------------------------------------------------------"
echo "Starting light control script..."
./home/pi/hand-hygiene/intervention/startup.sh

echo "--------------------------------------------------------------------------------"
exit 0
