# deploy/systemd/mqtt_imaging.service
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/home/pi/Desktop/picamera-mqtt/deploy/mqtt_imaging.sh
KillSignal=SIGINT
TimeoutStopSec=5
Restart=always
Type=forking
PIDFile=/var/run/mqtt_imaging.pid

[Install]
WantedBy=default.target
